# nomad-autoscaler-setup/ansible/roles/influxdb/templates/influxdb.service.j2
[Unit]
Description=InfluxDB - Time Series Database
Documentation=https://docs.influxdata.com/influxdb/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=influxdb
Group=influxdb
EnvironmentFile=-/etc/default/influxdb
ExecStart=/usr/local/bin/influxd -config {{ influxdb_config_dir }}/influxdb.conf
StandardOutput=journal
StandardError=journal
SyslogIdentifier=influxdb
Restart=on-failure
RestartSec=5s
KillMode=process

[Install]
WantedBy=multi-user.target
